# home/.config/systemd/user/awesome.service
[Unit]
Description=Awesome
After=hsetroot.service

[Service]
ExecStart=/usr/bin/awesome
Restart=always
RestartSec=10

[Install]
WantedBy=wm.target
