# home/.config/systemd/user/bitcoin-qt.service
[Unit]
Description=Bitcoin Wallet

[Service]
ExecStart=/usr/bin/bitcoin-qt -splash 0 -min

[Install]
WantedBy=mystuff.target
