# home/.config/systemd/user/conky.service
[Unit]
Description=Conky

[Service]
Type=forking
ExecStart=/usr/bin/conky
ExecStop=/usr/bin/killall conky

[Install]
WantedBy=mystuff.target
