# home/.config/systemd/user/nm-applet.service
[Unit]
Description=Network Manager applet

[Service]
ExecStart=/usr/bin/nm-applet

[Install]
WantedBy=mystuff.target
