# home/.config/systemd/user/skype.service
[Unit]
Description=Skype

[Service]
ExecStart=/usr/bin/skype

[Install]
WantedBy=mystuff.target
