# home/.config/systemd/user/transmission-gtk.service
[Unit]
Description=Transmission

[Service]
ExecStart=/usr/bin/transmission-gtk -m

[Install]
WantedBy=mystuff.target
