# home/.config/systemd/user/udiskie.service
[Unit]
Description=Handle automounting

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=mystuff.target
