# home/.config/systemd/user/wm.target
[Unit]
Description=Window manager target
AllowIsolate=true

[Install]
Alias=default.target
