# home/.config/systemd/user/xbindkeys.service
[Unit]
Description=bind hotkeys

[Service]
ExecStart=/usr/bin/xbindkeys -n -f %h/.xbindkeysrc

[Install]
WantedBy=mystuff.target
