# home/.config/systemd/user/xcompmgr.service
[Unit]
Description=Compositing manager

[Service]
ExecStart=/usr/bin/xcompmgr

[Install]
WantedBy=mystuff.target
