# home/.config/systemd/user/xrdb.service
[Unit]
Description=Xrdb

[Service]
Type=oneshot
ExecStart=/usr/bin/xrdb -merge .Xresources

[Install]
WantedBy=mystuff.target
